# docker/root/etc/s6-overlay/s6-rc.d/api/run
#!/usr/bin/with-contenv bash
exec /etc/s6-overlay/scripts/run-api
