# docker/root/etc/s6-overlay/s6-rc.d/nginx/run
#!/usr/bin/with-contenv bash
mkdir -p /config/log/nginx
exec nginx -g "daemon off;" 2>&1
